# openwrt-scripts/init.d/big-button
#!/bin/sh /etc/rc.common
# Big Internet Button - Init Script for OpenWrt
# Place in /etc/init.d/big-button

START=95
STOP=10
USE_PROCD=1

PROG=/usr/local/bin/big-button-daemon.sh
CONFIG_DIR=/etc/big-button
LOG_FILE=/tmp/big-button.log

# Create necessary directories and files
init_environment() {
    # Create config directory
    [ ! -d "$CONFIG_DIR" ] && mkdir -p "$CONFIG_DIR"
    
    # Create default config if not exists
    if [ ! -f "$CONFIG_DIR/config" ]; then
        cat > "$CONFIG_DIR/config" << EOF
# Big Internet Button Configuration
TIMER_MINUTES=40
WARNING_MINUTES=39
WARNING_REPEAT_SECONDS=30
DEVICE_SERIAL="/dev/ttyACM0"
DEVICE_INPUT="/dev/input/event0"
ENABLE_LOGGING=1
LOG_FILE="/tmp/big-button.log"
DEBOUNCE_TIME=2
EOF
    fi
    
    # Initialize state file - start blocked by default
    if [ ! -f "$CONFIG_DIR/state" ]; then
        # Load config to get TIMER_MINUTES value
        . "$CONFIG_DIR/config"
        echo "$TIMER_MINUTES" > "$CONFIG_DIR/state"
    fi
    [ ! -f "$CONFIG_DIR/state.status" ] && echo "blocked" > "$CONFIG_DIR/state.status"
    
    # Create log file
    [ ! -f "$LOG_FILE" ] && touch "$LOG_FILE"
}

start_service() {
    # Initialize environment
    init_environment
    
    # Load configuration
    . "$CONFIG_DIR/config"
    
    # Check if USB devices are present
    if [ ! -c "/dev/ttyACM0" ]; then
        echo "Warning: Serial device /dev/ttyACM0 not found"
        echo "Make sure the Big Internet Button is connected and drivers are loaded"
        
        # Try to load modules
        modprobe cdc_acm 2>/dev/null
        modprobe usbhid 2>/dev/null
        
        # Wait a moment for devices to appear
        sleep 2
        
        # Check again
        if [ ! -c "/dev/ttyACM0" ]; then
            echo "Error: Cannot find button device. Aborting."
            return 1
        fi
    fi
    
    # Start supervisor which will keep listener running
    /usr/local/bin/big-button-supervisor.sh > /tmp/supervisor.log 2>&1 &
    echo $! > /var/run/big-button-supervisor.pid
    
    # Setup cron job for timer
    grep -v 'big-button-timer' /etc/crontabs/root > /tmp/cron.tmp 2>/dev/null || true
    echo '* * * * * /usr/local/bin/big-button-timer.sh' >> /tmp/cron.tmp
    mv /tmp/cron.tmp /etc/crontabs/root
    /etc/init.d/cron restart >/dev/null 2>&1
    
    # Apply initial blocking if state is blocked
    if [ "$(cat $CONFIG_DIR/state.status 2>/dev/null)" = "blocked" ]; then
        echo "Applying initial internet blocking..."
        # Create nftables rules to block forwarding
        nft add table inet big_button 2>/dev/null
        nft add chain inet big_button forward \{ type filter hook forward priority 0 \; \} 2>/dev/null
        nft add rule inet big_button forward drop 2>/dev/null
        
        # Turn on red LED to indicate blocked state
        echo '2' > /dev/ttyACM0 2>/dev/null
        
        # Create blocked marker file
        touch "$CONFIG_DIR/state.blocked"
        
        echo "$(date '+%Y-%m-%d %H:%M:%S') - INIT: Internet blocked on startup" >> "$LOG_FILE"
    else
        # Give feedback - quick blink and beep
        echo '2' > /dev/ttyACM0 2>/dev/null
        sleep 1
        echo '1' > /dev/ttyACM0 2>/dev/null
        echo '3' > /dev/ttyACM0 2>/dev/null
    fi
    
    echo "Big Internet Button service started"
}

stop_service() {
    # Stop daemon
    if [ -f /var/run/big-button-daemon.pid ]; then
        kill $(cat /var/run/big-button-daemon.pid) 2>/dev/null
        rm -f /var/run/big-button-daemon.pid
    fi
    
    # Stop supervisor (which will stop listener too)
    if [ -f /var/run/big-button-supervisor.pid ]; then
        kill $(cat /var/run/big-button-supervisor.pid) 2>/dev/null
        rm -f /var/run/big-button-supervisor.pid
    fi
    
    # Also stop listener if running
    if [ -f /var/run/big-button-listener.pid ]; then
        kill $(cat /var/run/big-button-listener.pid) 2>/dev/null
        rm -f /var/run/big-button-listener.pid
    fi
    
    # Remove cron job
    grep -v "big-button-timer.sh" /etc/crontabs/root > /tmp/crontab.tmp 2>/dev/null
    mv /tmp/crontab.tmp /etc/crontabs/root 2>/dev/null
    /etc/init.d/cron restart
    
    # Ensure internet is unblocked
    nft delete table inet big_button 2>/dev/null
    
    # Turn off LED
    [ -c "/dev/ttyACM0" ] && echo "1" > /dev/ttyACM0 2>/dev/null
    
    echo "Big Internet Button service stopped"
}

restart() {
    stop
    sleep 1
    start
}

status() {
    if [ -f /var/run/big-button-daemon.pid ] && kill -0 $(cat /var/run/big-button-daemon.pid) 2>/dev/null; then
        echo "Big Internet Button is running"
        /usr/local/bin/big-button-control.sh status
    else
        echo "Big Internet Button is not running"
    fi
}

# Additional commands
boot() {
    # Wait for USB devices to be ready
    sleep 5
    start
}